# agent-go/packaging/deb/serverctl-agent.service
[Unit]
Description=ServerCtl Agent
After=network.target

[Service]
ExecStart=/usr/local/bin/serverctl-agent
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=serverctl-agent

[Install]
WantedBy=multi-user.target
